# systemd/qf-tv.service
# App starts from /home/kiosk/.config/openbox/autostart after LightDM creates :0.
# Do not enable this unit to ExecStart qf_tv — that races and logs "cannot open display: :0".
# Recovery: sudo bash scripts/fix-tv-display.sh

[Unit]
Description=QueueFlow TV (placeholder — openbox runs qf_tv)
After=lightdm.service
PartOf=lightdm.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=graphical.target
